# RemAssist/soa1-api.service
# SOA1 API Service - Systemd Unit File
# 
# Installation:
#   sudo cp RemAssist/soa1-api.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soa1-api
#   sudo systemctl start soa1-api
#
# Logs:
#   journalctl -u soa1-api -f

[Unit]
Description=SOA1 Home Assistant API Service
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=ryzen
Group=ryzen
WorkingDirectory=/home/ryzen/projects/home-ai/soa1
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ryzen/projects:/home/ryzen/projects/home-ai
ExecStart=/usr/bin/env uvicorn api:create_app --workers 1 --host 0.0.0.0 --port 8001 --factory --log-level info
Restart=on-failure
RestartSec=5

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ryzen/projects/home-ai/finance-agent/data
ReadWritePaths=/home/ryzen/projects/home-ai/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
